# scholarshipowl/docs/sowl.localtunnel.conf
[program:sowlpavelz.localtunnel]
command=lt -s sowlpavelz -p 80
stdout_logfile=/home/vagrant/scholarship/storage/localtunnel.log
redirect_stderr=true
autostart=true
autorestart=true
